# Library/LaunchAgents/com.dotfiles.ollama-service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Optional warm Ollama daemon for local model workflows. -->

    <key>Label</key>
    <string>com.dotfiles.ollama-service</string>

    <key>ProgramArguments</key>
    <array>
        <string>__DOTFILES_ROOT__/scripts/services/ollama-service.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__HOME__</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>__HOME__/.local/state/ollama/service.out.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.local/state/ollama/service.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:__HOME__/.local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>OLLAMA_HOST</key>
        <string>127.0.0.1:11434</string>
    </dict>
</dict>
</plist>
