# Library/LaunchAgents/com.kubectl-fzf-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kubectl-fzf-server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>export PATH="$HOME/go/bin:$PATH"; kubectl-fzf-server</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>StandardOutPath</key>
    <string>/tmp/kubectl-fzf-server.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/kubectl-fzf-server.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/shaheislam</string>
        <key>KUBECONFIG</key>
        <string>/Users/shaheislam/.kube/config</string>
        <key>KUBECTL_FZF_CACHE</key>
        <string>/tmp/kubectl_fzf_cache</string>
    </dict>
</dict>
</plist>
