# Library/LaunchAgents/com.dotfiles.pinchtab-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Persistent PinchTab orchestrator for browser automation across agents. -->
    <!-- Install: launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.dotfiles.pinchtab-server.plist -->
    <!-- Restart: launchctl kickstart -k gui/$(id -u)/com.dotfiles.pinchtab-server -->

    <key>Label</key>
    <string>com.dotfiles.pinchtab-server</string>

    <key>ProgramArguments</key>
    <array>
        <string>__DOTFILES_ROOT__/scripts/pinchtab/server.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__HOME__</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>__HOME__/.local/state/pinchtab/server.out.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.local/state/pinchtab/server.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:__HOME__/.bun/bin:__HOME__/.local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PINCHTAB_PORT</key>
        <string>9867</string>
    </dict>
</dict>
</plist>
