# Library/LaunchAgents/com.dotfiles.changelog-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Monthly changelog review: queues a gwt-ticket for changelog analysis -->
    <!-- Runs on the 1st of each month at 10:00 local time -->
    <!-- Install: launchctl load ~/Library/LaunchAgents/com.dotfiles.changelog-review.plist -->
    <!-- Remove:  launchctl unload ~/Library/LaunchAgents/com.dotfiles.changelog-review.plist -->

    <key>Label</key>
    <string>com.dotfiles.changelog-review</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/fish</string>
        <string>-c</string>
        <string>gwt-queue add TASK "Monthly changelog review" "Review tool changelogs from the last 30 days, update dotfiles config where warranted, and create atomic commits per tool." --template changelog-review --no-beads</string>
    </array>

    <!-- 1st of each month at 10:00 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/changelog-review-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/changelog-review-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
		<key>HOME</key>
		<string>__HOME__</string>
	</dict>
</dict>
</plist>
